# rootfs/etc/s6-overlay/s6-rc.d/init-ffmpeg/run
#!/command/with-contenv bashio
# shellcheck shell=bash

eval "$(bashio::config 'command')"


